# scripts/deployment/club-corra-api.service
[Unit]
Description=Club Corra API Service
After=network.target postgresql.service
Wants=postgresql.service
StartLimitIntervalSec=0

[Service]
Type=simple
User=ec2-user
Group=ec2-user
WorkingDirectory=/opt/club-corra-api
ExecStart=/usr/bin/node dist/main.js
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StartLimitBurst=5
StartLimitInterval=60

# Environment variables from production config
Environment=NODE_ENV=production
Environment=PORT=8080
Environment=HOST=127.0.0.1
Environment=SSL_CERT_PATH=/etc/letsencrypt/live/16.170.179.71.nip.io/fullchain.pem
Environment=SSL_KEY_PATH=/etc/letsencrypt/live/16.170.179.71.nip.io/privkey.pem
Environment=CORS_ORIGIN=https://club-corra-pilot-admin.vercel.app,https://club-corra-pilot-admin-git-master-vikas-ahlawats-projects.vercel.app,https://club-corra-pilot-admin-d8hqxkw2x-vikas-ahlawats-projects.vercel.app,https://*.vercel.app,https://admin.clubcorra.com,https://clubcorra.com,https://*.clubcorra.com

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=club-corra-api

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/club-corra-api

[Install]
WantedBy=multi-user.target
